# supervisord/troup-nodes.cfg
[program:troup-node]
command=python3 -m troup --node troup-node-%(host_node_name)s-%(process_num)s --port 700%(process_num)d --lock
process_name=troup-node
autostart=true
autorestart=true
user=troup
directory=/var/lib/troup/node
stdout_logfile=/var/log/troup.node.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=100
stderr_logfile=/var/log/trou.node.err.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=100
